# .config/systemd/user/mbsync.service
[Unit]
Description=Run mbsync to sync mail

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -a
